# infrastructure/depth_node/systemd/janus-camera-page.service.d/override.conf
[Unit]
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Restart=on-failure
RestartSec=5
LimitNOFILE=4096
MemoryMax=512M
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
ReadWritePaths=/var/lib/camera-fdir
ReadWritePaths=/run/camera
