# infrastructure/color_node/systemd/janus-camera-page.service.d/override.conf
[Unit]
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Environment=PYTHONPATH=/opt/janus-camera-page
# Enforce operator's "no reboots" constraint: FDIR level 4 (reboot_node)
# bypasses physical reboot, transitioning the system to SAFE mode instead.
# Without this var, recovery_ladder could physically invoke systemctl reboot
# when L0-L3 escalation is exhausted. See janus_camera_page/CONTRACT.md.
Environment=CAM_WATCHDOG_REBOOT_ENABLED=0
# Explicit default instance for encoder-admin CLI (avoid silent hardcoded
# fallback). Single-camera deployment = cam-rgb. Multi-camera: override
# per-call via --instance arg.
Environment=ENCODER_DEFAULT_INSTANCE=cam-rgb
# Sprint X3.4: stack default = joystick off. This robot wraps stack with joystick
# enabled (we ARE the robot deploy — not a generic stack consumer). Override
# to "always" keeps existing UX. Generic stack users leave unset = "off".
Environment=STACK_DEFAULT_JOYSTICK_MODE=always
# Service-layer LAN invariants (review P0-4): the gateway's own LAN IP and the
# camera subnet. add_node_by_host rejects the gateway IP and any host outside this
# CIDR — enforced in the store so route/fleet/future callers are all bound.
Environment=GATEWAY_LAN_IP=192.168.1.10
Environment=CAMERA_LAN_CIDR=192.168.1.0/24
EnvironmentFile=-/etc/robot/camera-secrets.env
# Track A: runtime-tunable NON-secret knobs (ICE_POLICY, TURN_CRED_TTL). Writable
# (/etc/robot is RW below); read live by Settings via default_factory + cache_clear.
# MUST NOT be re-added as `Environment=ICE_POLICY=` above — a directive shadows this file.
EnvironmentFile=-/etc/robot/rs-runtime.env
ExecStart=
ExecStart=/opt/janus-camera-page/.venv/bin/python /opt/janus-camera-page/main.py
StateDirectory=camera-fdir
Restart=on-failure
RestartSec=5

# Resource limits
LimitNOFILE=4096
MemoryMax=512M

# Sandboxing
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=yes
ReadWritePaths=/var/lib/camera-fdir
ReadWritePaths=/run/camera
# ProtectSystem=strict makes the whole FS (incl. /tmp) read-only, so Python's
# tempfile (used by ssh_transport for the pinned known_hosts during node
# provision/activate) has no writable temp dir → "No usable temporary directory".
# Point TMPDIR at the writable runtime dir. (/tmp stays readable for the node
# bundle at /tmp/camera-node-bundle.tar.gz; PrivateTmp would hide that, so not it.)
Environment=TMPDIR=/run/camera
# FDIR event persistence (fdir_events._persist) writes /var/log/camera-fdir/fdir.jsonl.
# Without this ReadWritePath, ProtectSystem=strict silently blocks writes → events lost.
# Symptom: /var/log/camera-fdir/fdir.jsonl mtime frozen, journal OK but history gone.
ReadWritePaths=/var/log/camera-fdir
# Phase 2: audit log path (structured admin action audit trail).
ReadWritePaths=/var/log/camera-audit
# /etc/robot/cam-rgb.tuning.env is written by L4 config endpoint (camera-page
# UI: /api/v1/{cam_type}/config POST). env_store.write_env_atomic creates a
# tempfile in parent dir then atomic-renames — covers tuning.env, the flock
# file (.cam-rgb.env.lock) and transient cam-rgb.*.env. Without this, /config
# POST fails with EROFS (sandbox-level RO, not disk).
ReadWritePaths=/etc/robot
# Defense-in-depth: /etc/robot is now RW (for tuning.env), but camera-secrets.env
# must stay RO so a compromised L4 process can't rewrite admin token / TURN creds.
BindReadOnlyPaths=/etc/robot/camera-secrets.env
